# ROOTFS/etc/systemd/system/rpcbind.service.d/override.conf
[Unit]
# ensure service doesn't run unless the firewall is up
After=waggle-firewall.service
BindsTo=waggle-firewall.service
